# templates/influxdb-systemd.service.j2
# ============================================================================
# InfluxDB Systemd Service Unit File
# ============================================================================
# ⚠️ WARNING: This file is managed by Ansible - DO NOT EDIT MANUALLY!       ⚠️
# Any manual changes will be overwritten on the next Ansible run.
# Edit the template file or Ansible variables instead.
#
# Template path: roles/databases/templates/influxdb-systemd.service.j2
# Deploy path: /etc/systemd/system/{{ influxdb_systemd_name }}.service
# Last deployed: {{ ansible_date_time.iso8601 }}
# ============================================================================
#
# This systemd service file configures InfluxDB to run as a system service,
# providing a time-series database for metrics and events.
#
# Features:
#   - Fast, scalable time-series storage
#   - HTTP API for data ingestion and queries
#
# Configuration:
#   - Config file: {{ influxdb_config_file }}
#   - Data directory: {{ influxdb_data_dir }}
#   - User: {{ influxdb_user }}
#   - Working directory: {{ influxdb_data_dir }}
#   - Binary location: {{ influxdb_binary_path }}
#   - Logs: sudo journalctl -u {{ influxdb_systemd_name }} -f
#
# Usage:
#   sudo systemctl enable {{ influxdb_systemd_name }}    # Enable on boot
#   sudo systemctl start {{ influxdb_systemd_name }}     # Start service
#   sudo systemctl status {{ influxdb_systemd_name }}    # Check status
#   sudo systemctl stop {{ influxdb_systemd_name }}      # Stop service
#   sudo journalctl -u {{ influxdb_systemd_name }}       # View logs
#   sudo journalctl -u {{ influxdb_systemd_name }} -f    # View and follow logs
#
# Security Notes:
#   - InfluxDB runs with limited privileges
#
# Dependencies:
#   - Requires network.target
#
# References:
#   - https://docs.influxdata.com/influxdb/
#
# Generated by Ansible on {{ ansible_date_time.iso8601 }}
# Host: {{ ansible_fqdn }}
# ------------------------------------------------------------------------------

[Unit]
Description=InfluxDB Time Series Database
Documentation=https://docs.influxdata.com/influxdb/
After=network.target
Requires=network.target

[Service]

# Set config file path via environment variable (for InfluxDB 2.x/3.x)
Environment=INFLUX_CONFIG_PATH={{ influxdb_config_file }}
{# Environment=INFLUXDB3_OBJECT_STORE=file
Environment=INFLUXDB3_OBJECT_STORE_PATH={{ influxdb_data_dir }}/object_store #}

ExecStart={{ influxdb_binary_path }} serve \
                                     --node-id {{ influxdb_systemd_name }}-01 \
                                     --data-dir={{ influxdb_data_dir }} \
                                     --object-store=file \
                                     --http-bind={{ databases_influxdb.listen_address }}:{{ databases_influxdb.port }}

Type=simple
Restart=on-failure
RestartSec=10
User={{ influxdb_user }}
Group={{ influxdb_group }}

# Directory creation and permissions
RuntimeDirectory={{ influxdb_systemd_name }}
RuntimeDirectoryMode=0750
ReadOnlyPaths={{ influxdb_config_dir }}
ReadWritePaths={{ influxdb_data_dir }}

# Logging to systemd
StandardOutput=journal
StandardError=journal
SyslogIdentifier={{ influxdb_systemd_name }}

# Hardening measures
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
NoNewPrivileges=true
PrivateDevices=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET
RestrictRealtime=true

[Install]
WantedBy=multi-user.target