# templates/mariadb-systemd.service.j2
# ============================================================================
# MariaDB Database Server Systemd Service Unit File
# ============================================================================
# ⚠️ WARNING: This file is managed by Ansible - DO NOT EDIT MANUALLY!       ⚠️
# Any manual changes will be overwritten on the next Ansible run.
# Edit the template file or Ansible variables instead.
#
# file_path: roles/databases/templates/mariadb-systemd.service.j2
# =============================================================================
#
# This systemd service file configures MariaDB to run as a system service,
# providing a robust, open source relational database system.
#
# Features:
#   - ACID compliant transactions
#   - MySQL compatibility
#   - High performance and scalability
#
# Configuration:
#   - Config file: {{ mariadb_config_dir }}/my.cnf
#   - Data directory: {{ mariadb_data_dir }}
#   - User: {{ mariadb_user }}
#   - Working directory: {{ mariadb_data_dir }}
#   - Binary location: /usr/local/bin/mysqld
#   - Logs: sudo journalctl -u {{ mariadb_systemd_name }} -f
#
# Usage:
#   sudo systemctl enable {{ mariadb_systemd_name }}    # Enable on boot
#   sudo systemctl start {{ mariadb_systemd_name }}     # Start service
#   sudo systemctl status {{ mariadb_systemd_name }}    # Check status
#   sudo systemctl stop {{ mariadb_systemd_name }}      # Stop service
#   sudo journalctl -u {{ mariadb_systemd_name }}       # View logs
#   sudo journalctl -u {{ mariadb_systemd_name }} -f    # View and follow logs
#
# Security Notes:
#   - MariaDB runs with limited privileges
#
# Dependencies:
#   - Requires network.target
#
# References:
#   - https://mariadb.com/kb/en/systemd/
#   - https://mariadb.com/kb/en/configuring-mariadb-with-option-files/
#
# Generated by Ansible on {{ ansible_date_time.iso8601 }}
# Host: {{ ansible_fqdn }}
# ------------------------------------------------------------------------------

[Unit]
Description=MariaDB Database Server
Documentation=https://mariadb.com/kb/en/

# Service Dependencies
After=network.target
Requires=network.target

[Service]
WorkingDirectory={{ mariadb_data_dir }}
ExecStart=/usr/local/bin/mysqld --defaults-file={{ mariadb_config_dir }}/{{ mariadb_user }}.conf
Type=simple
Restart=on-failure
RestartSec=10

# Directory creation and permissions
User={{ mariadb_user }}
Group={{ mariadb_user }}
RuntimeDirectory={{ mariadb_user }}
RuntimeDirectoryMode=0750
ReadOnlyPaths={{ mariadb_config_dir }}
ReadWritePaths={{ mariadb_data_dir }} {{ mariadb_run_dir }}

# Logging to systemd
StandardOutput=journal
StandardError=journal
SyslogIdentifier={{ mariadb_user }}

# Hardening measures
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
NoNewPrivileges=true
PrivateDevices=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET
RestrictRealtime=true

[Install]
WantedBy=multi-user.target
