# templates/mongodb-systemd.service.j2
# ============================================================================
# MongoDB systemd Service Unit (Ansible managed)
# ============================================================================
# ⚠️ WARNING: This file is managed by Ansible - DO NOT EDIT MANUALLY!       ⚠️
# Any manual changes will be overwritten on the next Ansible run.
# Edit the template file or Ansible variables instead.
#
# Template path: roles/databases/templates/mongodb-systemd.service.j2
# Deploy path: /lib/systemd/system/{{ mongodb_systemd_name }}.service
# Last deployed: {{ ansible_date_time.iso8601 }}
# =============================================================================
#
# This systemd service file configures MongoDB to run as a system service.
#
# Features:
#   - Runs mongod as a dedicated system user
#   - Logs to systemd journal
#   - Hardened with systemd security options
#
# Configuration:
#   - Config file: {{ mongodb_config_file }}
#   - Data directory: {{ mongodb_data_dir }}
#   - User: {{ mongodb_user }}
#   - Binary: /usr/local/bin/mongod
#   - Logs: sudo journalctl -u {{ mongodb_systemd_name }} -f
#
# Usage:
#   sudo systemctl enable {{ mongodb_systemd_name }}
#   sudo systemctl start {{ mongodb_systemd_name }}
#   sudo systemctl status {{ mongodb_systemd_name }}
#   sudo systemctl stop {{ mongodb_systemd_name }}
#   sudo journalctl -u {{ mongodb_systemd_name }} -f
#
# Security Notes:
#   - mongod runs with limited privileges
#
# Dependencies:
#   - Requires network.target
#
# References:
#   - https://docs.mongodb.com/manual/tutorial/manage-mongodb-processes/
#
# Generated by Ansible on {{ ansible_date_time.iso8601 }}
# Host: {{ ansible_fqdn }}
# ------------------------------------------------------------------------------

[Unit]
Description=MongoDB Database Server
Documentation=https://docs.mongodb.com/manual/
After=network.target
Requires=network.target

[Service]
ExecStart=/usr/local/bin/mongod --config {{ mongodb_config_file }} \
                                --dbpath {{ mongodb_data_dir }} \
                                # --pidfilepath /run/{{ mongodb_systemd_name }}/mongod.pid

# Process management
####################
Type=simple
Restart=on-failure
RestartSec=10
User={{ mongodb_user }}
Group={{ mongodb_group }}


# Directory creation and permissions
####################################

RuntimeDirectory={{ mongodb_systemd_name }}
RuntimeDirectoryMode=0750
ReadOnlyPaths={{ mongodb_config_dir }}
ReadWritePaths={{ mongodb_data_dir }}


# Logging to systemd
####################
StandardOutput=journal
StandardError=journal
SyslogIdentifier={{ mongodb_systemd_name }}


# Hardening measures
####################
# PrivateTmp=true
# ProtectSystem=full
# ProtectHome=true
# NoNewPrivileges=true
# PrivateDevices=true
# MemoryDenyWriteExecute=true
# SystemCallArchitectures=native
# ProtectKernelModules=true
# ProtectControlGroups=true
# RestrictAddressFamilies=AF_INET
# RestrictRealtime=true

[Install]
WantedBy=multi-user.target
